# post-startup.service
[Unit]
Description=Provision /var state before services start
After=local-fs.target network.target
Before=postgresql.service caddy.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/post-startup.sh
RemainAfterExit=yes

ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/lib/pgsql /var/lib/webdav /etc/caddy /run/postgresql

[Install]
WantedBy=multi-user.target
