# post-startup.service
[Unit]
Description=Provision /var state before services start
After=local-fs.target
Before=httpd.service postgresql.service caddy.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/post-startup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
