# post-startup-root.service
[Unit]
Description=Root-side boot setup (caddy env, postgres hash)
Requires=network-online.target
After=local-fs.target network-online.target nftables.service
Before=post-startup-postgresql.service caddy.service

[Service]
Type=oneshot
User=root
Group=creds
ExecStart=/opt/scripts/post-startup-root.sh
RemainAfterExit=yes
RuntimeDirectory=post-startup-postgresql
RuntimeDirectoryMode=0750

ProtectSystem=strict
ProtectHome=true
ProtectProc=invisible
ProtectClock=true
ProtectKernelLogs=true
ProtectKernelTunables=true
ProtectHostname=true
PrivateDevices=true
PrivateTmp=true
PrivateIPC=true
RestrictNamespaces=true
NoNewPrivileges=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET
SystemCallFilter=@system-service
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
ReadWritePaths=/run/caddy /etc/rclone
ReadOnlyPaths=-/var/lib/pgsql/data/base

[Install]
WantedBy=multi-user.target
